# service/arch-os-update.timer
[Unit]
Description=Run Arch OS update check every hour

[Timer]
OnBootSec=30sec
OnUnitActiveSec=120min

[Install]
WantedBy=timers.target